# docker/cron/root.crontab
SHELL=/bin/sh
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
APP_DIR=/var/www/app

*/30 * * * * /usr/local/bin/run-console-command.sh income_rates_refresh app:income:rates:refresh-live >> /proc/1/fd/1 2>> /proc/1/fd/2
2,32 * * * * /usr/local/bin/run-console-command.sh income_live_fill app:income:fill-live-gel >> /proc/1/fd/1 2>> /proc/1/fd/2
17 * * * * /usr/local/bin/run-console-command.sh income_official_backfill app:income:backfill-gel >> /proc/1/fd/1 2>> /proc/1/fd/2
